# canbus.service
[Unit]
Description=Setup CAN interface and run CAN tranceiver script
After=network.target

[Service]
ExecStart=/home/USERNAME/pi-canbus-transceiver/setup.sh
WorkingDirectory=/home/USERNAME/pi-canbus-transceiver
StandardOutput=inherit
StandardError=inherit
Restart=always
User=USERNAME

[Install]
WantedBy=multi-user.target